# config/supervisord.ini
[program:gunicorn]
command=/usr/local/bin/gunicorn --bind 0.0.0.0:8081 -w 1 -t 5040 ums.__main__:main() --access-logfile /var/log/supervisord/ums.log
directory=/usr/src/app/
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisord/ums.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisord/ums.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
